# vessence/vault_web/vault-web.service
# NOTE: This is a reference template only. Paths and User must be customized for deployment.
[Unit]
Description=Amber Vault Web Server
After=network.target

[Service]
Type=simple
User=chieh
WorkingDirectory=/home/chieh/vessence/vault_web
ExecStart=/home/chieh/google-adk-env/adk-venv/bin/uvicorn main:app --host 127.0.0.1 --port 8080 --proxy-headers --forwarded-allow-ips='*'
Restart=on-failure
RestartSec=5
StandardOutput=append:/home/chieh/vessence/logs/vault_web.log
StandardError=append:/home/chieh/vessence/logs/vault_web.log
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
